# service/tail/run
#!/bin/sh
tail -f /etc/service/*/log/current
